# dw/config/supervisord.conf
[supervisord]
nodaemon=true

[program:dw]
directory=/services
command= java -jar dropwizard-example-0.9.0-rc2-SNAPSHOT.jar server example.yml
