# kibana/config/supervisord.conf
[supervisord]
nodaemon=true

[program:kibana]
command=/kibana-4.1.0-linux-x64/bin/kibana
autorestart=true

[program:nginx]
command=/usr/sbin/nginx
autorestart=true